# modular-agents.service
[Unit]
Description=Modular Agents — Semi Autonomous AI Assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/home/YOUR_USERNAME/modular-agents
Environment="PATH=/home/YOUR_USERNAME/modular-agents/.venv/bin"
ExecStart=/home/YOUR_USERNAME/modular-agents/.venv/bin/python /home/YOUR_USERNAME/modular-agents/main.py
Restart=on-failure
RestartSec=10
StartLimitInterval=120
StartLimitBurst=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
