# modular-agents-backup.timer
[Unit]
Description=Modular Agents — daily backup timer
Requires=modular-agents-backup.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
